# systemd/selfhostbox.service
[Unit]
Description=SelfHostBox - Self-Hosting Platform
After=network.target docker.service

[Service]
Type=simple
User=selfhostbox
Group=selfhostbox
WorkingDirectory=/opt/selfhostbox
Environment=PYTHONPATH=/opt/selfhostbox
Environment=SELFHOSTBOX_CONFIG=/opt/selfhostbox/config/config.yml
ExecStart=/opt/selfhostbox/.venv/bin/gunicorn --bind 0.0.0.0:8080 --workers 2 "backend.app:create_app()"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target